# aegis.service
[Unit]
Description=Aegis Gateway — embedded cluster mode (Gateway + Mesh Relay + TRUSTMARK Engine + Botawiki)
After=network.target nats.service

[Service]
Type=simple
User=aegis
ExecStart=/usr/local/bin/aegis-gateway -c /etc/aegis/gateway.toml --embedded
Restart=on-failure
RestartSec=3
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
